# deploy/systemd/ailab-mcp.service
[Unit]
Description=ailab MCP Server (HTTP transport for Hermes Agent)
Documentation=file:///home/admin/KitchenBook/backend/mcp_server/README.md
After=network.target gunicorn.service
Wants=gunicorn.service

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/KitchenBook/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings
Environment=PATH=/home/admin/KitchenBook/venv/bin:/usr/bin:/bin
EnvironmentFile=-/home/admin/KitchenBook/backend/.env
# 用 fastmcp CLI 起 server —— fastmcp 3.x 的 mcp.run() 不再接受 host/port kwarg，
# CLI 反而原生支持 --transport/--host/--port。0 侵入项目代码。
#
# Bind 0.0.0.0 是为了让 docker 容器（从 bridge 网关 IP 172.18.0.1 访问）
# 能连到 MCP。公网入站靠阿里云安全组 + iptables 拦截 8100 端口。
ExecStart=/home/admin/KitchenBook/venv/bin/fastmcp run mcp_server/server.py:mcp \
          --transport http --host 0.0.0.0 --port 8100 --no-banner
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# ===== 收紧 =====
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ailab-mcp

[Install]
WantedBy=multi-user.target
